# .crontab
* * * * * /Users/mminakova/Documents/COVID/covid_tracker/covid/bash.sh
